# deploy/skysaver.service
[Unit]
Description=SkySaver Flight Price AI Agent — FastAPI Service
Documentation=https://github.com/your-repo/flight-agent
After=network.target
Wants=network-online.target

[Service]
Type=exec
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/flight-agent
Environment="PATH=/home/ubuntu/flight-agent/venv/bin"
EnvironmentFile=/home/ubuntu/flight-agent/.env

ExecStart=/home/ubuntu/flight-agent/venv/bin/gunicorn \
    -c gunicorn_conf.py \
    "api.main:create_app()"

# Restart policy — restart on any failure, wait 10s between restarts.
# StartLimitBurst=3 prevents infinite restart loops in < 60s windows.
Restart=always
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=3

# Resource limits — prevent OOM from taking down the whole VPS.
# 24 GB total: 8GB for SkySaver, leaving 16GB for n8n + OS + headroom.
MemoryMax=8G
# Allow up to 3 of 4 vCPUs (leave 1 for n8n + OS scheduler).
CPUQuota=300%

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/ubuntu/flight-agent

# Logging — captured by journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=skysaver

[Install]
WantedBy=multi-user.target
